# daily.crontab
0 21 * * * cd ~/prjs/job_hooroom && ./etl.sh
